# hive-update.service
[Unit]
Description=Hive Auto Update Service
After=network.target evaluasi.service
Documentation=https://github.com/barangbaru/solid-apps

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hive-update-run.sh
StandardOutput=append:/tmp/hive_update.log
StandardError=append:/tmp/hive_update.log
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
